# supervisor/ltc2983-ioc-02.conf
[program:ltc2983-ioc-02]
directory=/opt/LTC2983-ioc/iocBoot/iocLTC2983
command=/opt/LTC2983-ioc/bin/linux-arm/LTC2983 ioc-spi02.cmd
stdout_logfile=/var/log/ltc2983-ioc-02.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stderr_logfile=/var/log/ltc2983-ioc-02.log
stderr_ogfile_maxbytes=1MB
stderr_logfile_backups=5
